# etc/system/etc/cron.d/update-ip
*/5 * * * * eliot curl ipecho.net/plain > /home/eliot/Dropbox/tadpole-ip.txt
